# news-crawler/crontab
# Run at 9 AM KST (Docker container uses UTC)
* * * * * echo `date` >> /logs/test.log
0 0,4,8 * * * /usr/local/bin/python3 /news-crawler/crawler.py >> /logs/news.log
